# goclaw.service
[Unit]
Description=GoClaw Gateway Service
After=network.target

[Service]
Type=simple
User=goclaw
Group=goclaw
# Adjust path to binary
ExecStart=/usr/local/bin/goclaw
Restart=always
RestartSec=5s

# Environment variables
Environment=GOCLAW_PORT=18789
Environment=GOCLAW_BIND=loopback
# Environment=GOCLAW_TOKEN=your-secret-token
# Environment=DISCORD_BOT_TOKEN=your-discord-token
# Environment=DISCORD_GUILD_ID=your-guild-id

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
